# deploy/api-data-process.service
[Unit]
Description=API Data Process Service (Gunicorn + Uvicorn)
After=network.target redis.service mysql.service
Wants=redis.service
StartLimitIntervalSec=60s
StartLimitBurst=3

[Service]
Type=notify
NotifyAccess=main

# ── Paths ─────────────────────────────────────────────────────────────────────
WorkingDirectory=/opt/api-data-process

# ── Environment ───────────────────────────────────────────────────────────────
EnvironmentFile=/opt/api-data-process/.env
Environment=PYTHONPATH=/opt/api-data-process/src
Environment=PYTHONUNBUFFERED=1

# ── Executable ────────────────────────────────────────────────────────────────
ExecStart=/opt/api-data-process/.venv/bin/gunicorn \
    --config /opt/api-data-process/gunicorn.conf.py \
    main:app

ExecReload=/bin/kill -s HUP $MAINPID

# ── Restart policy ────────────────────────────────────────────────────────────
Restart=on-failure
RestartSec=5s

# ── Resource limits ───────────────────────────────────────────────────────────
LimitNOFILE=65535
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
